# debian/postinst
#!/bin/bash
#
# postinst script for linuxmuster-nginx-proxy
# thomas@linuxmuster.net
# 20180312
# GPL v3
#

# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

name="linuxmuster-nginx-proxy"
images="jwilder/nginx-proxy jwilder/whoami"
sharedir="/usr/share/$name"
dockerkey="/etc/linuxmuster/ssl/docker.key.pem"
dockerdir="/srv/docker/$name"

case "$1" in

 configure)

  # pull docker image
  if ping -q -c2 google.com &> /dev/null; then
    for i in $images; do
      docker pull "$i"
    done
  fi

  if [ ! -e "$dockerkey" ]; then
    echo "SSL hostkey not found!"
    echo "Please run linuxmuster-prepare with -c option first,"
    echo "then run dpkg-reconfigure $name."
    exit 0
  fi

  # create docker group & user
  grep -q ^docker /etc/group || groupadd -g 5000 docker
  grep ^docker /etc/group | grep -q 5000 || groupmod -g 5000 docker
  grep -q ^docker /etc/passwd || useradd -u 5000 -d /srv/docker -g docker docker
  grep ^docker /etc/passwd | grep -q 5000 || usermod -u 5000 -d /srv/docker -g docker docker

  # enable and start service
  /bin/systemctl enable "$name.service"
  /bin/systemctl stop "$name.service"
  /bin/systemctl start "$name.service"

 ;;

 abort-upgrade|abort-remove|abort-deconfigure)
 ;;

 *)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
 ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
